# config/supervisor/mux.conf
[program:Multiplex]
command=odr-dabmux /home/azerty/dab/config/conf.mux
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/mux.log
stdout_logfile=/var/log/supervisor/mux.log
